# deploy/zrws.service
[Unit]
Description=ZRWS 智壤卫士审批服务
After=network.target mysql.service redis.service

[Service]
Type=simple
User=root
WorkingDirectory=/app
ExecStart=/usr/bin/java -Xms512m -Xmx1024m -jar /app/app.jar
ExecStop=/bin/kill -15 $MAINPID
Restart=on-failure
RestartSec=10
StandardOutput=file:/app/logs/zrws.log
StandardError=file:/app/logs/zrws-error.log

# 环境变量
Environment="JAVA_HOME=/usr/lib/jvm/java-17-openjdk"
Environment="SPRING_PROFILES_ACTIVE=prod"
Environment="OPENAI_API_KEY=your-api-key"
Environment="REDIS_HOST=localhost"
Environment="REDIS_PORT=6379"

[Install]
WantedBy=multi-user.target